# Projeto8-LXC/containers/lxc-fedora-http/rootfs/etc/systemd/system/configure-network.service
[Unit]
Description=Configure Network at Boot
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/configure-network.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
